# docker-config/php/crontab/cron
* * * * * /var/www/docker-config/php/shell/notify.sh
* * * * * /var/www/docker-config/php/shell/test.sh
