# docker-config/php/supervisor/cron-worker.conf
[program:cron-worker]

command=/bin/bash /var/www/docker-config/php/shell/setup.sh

autostart=true

autorestart=true

user=root

log_stderr=true

logfile=/var/log/testcron.log
